# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fastapi]
command=uvicorn main:app --host 0.0.0.0 --port 8056
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/var/log/supervisor/fastapi_error.log
stdout_logfile=/var/log/supervisor/fastapi_out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
redirect_stderr=false

[program:worker]
command=python -m utils.worker
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/var/log/supervisor/worker_error.log
stdout_logfile=/var/log/supervisor/worker_out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
redirect_stderr=false